# wamp/crossbar.ini
[program:crossbar]
command=crossbar start --cbdir /var/www/code/.crossbar
autostart=true
autorestart=true
